# dot_config/autostart/private_ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/id_rsa
Name=ssh-add
Type=Application
