# dot_config/autostart/private_ssh-add.desktop
[Desktop Entry]
Exec=ssh-add -q ~/.ssh/id_rsa ~/.ssh/id_rsa_deploy
Name=ssh-add
Type=Application
